# kebbi_flask.service
[Unit]
Description=Kebbi Emotion Flask GPU Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/ubuntu/server_gpu.py
WorkingDirectory=/home/ubuntu/
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
